# http_role/templates/kea-status.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Refresh the PXE/DHCP status page

[Timer]
OnBootSec=1min
OnUnitActiveSec={{ dashboard_refresh }}

[Install]
WantedBy=timers.target
